# etc/xdg/autostart/tuquito-xhost-plus.desktop.desktop
[Desktop Entry]
Encoding=UTF-8
Version=1.0
Name=Xhost +
Exec=xhost +
Terminal=false
Type=Application
StartupNotify=false
Terminal=false
